# docker/supervisor/supervisord.conf
[supervisord]
logfile = /var/www/storage/logs/supervisord.log
logfile_maxbytes = 50MB
logfile_backups=10
loglevel = info
pidfile = /tmp/supervisord.pid
nodaemon = true
minfds = 1024
minprocs = 200
umask = 022
user = root
identifier = supervisor
directory = /tmp
nocleanup = true
childlogdir = /tmp
strip_ansi = false

[program:explorer.worker.blocks]
command=php artisan minter:api:pull-node-data -vvv
process_name=%(program_name)s_%(process_num)02d
numprocs=1
directory=/var/www
autostart=true
autorestart=true
startretries=50
user=root
stdout_logfile=/var/www/storage/logs/blocks.out.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/www/storage/logs/blocks.err.log
stderr_logfile_maxbytes=10MB

[program:explorer.worker.nodes]
command=php artisan minter:check-node-list -vvv
process_name=%(program_name)s_%(process_num)02d
numprocs=1
directory=/var/www
autostart=true
autorestart=true
user=root
stdout_logfile=/var/www/storage/logs/node.out.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/www/storage/logs/node.err.log
stderr_logfile_maxbytes=10MB

[program:explorer.worker.transactions]
command=php artisan queue:work --queue=transactions --tries=3 --daemon
process_name=%(program_name)s_%(process_num)02d
directory=/var/www
autostart=true
autorestart=true
user=root
numprocs=10
stdout_logfile=/var/www/storage/logs/jobs.out.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/www/storage/logs/jobs.err.log
stderr_logfile_maxbytes=10MB

[program:explorer.worker.balance]
command=php artisan queue:work --queue=balance --tries=3 --daemon
process_name=%(program_name)s_%(process_num)02d
directory=/var/www
autostart=true
autorestart=true
user=root
numprocs=5
stdout_logfile=/var/www/storage/logs/jobs.out.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/www/storage/logs/jobs.err.log
stderr_logfile_maxbytes=10MB

[program:explorer.worker.validators]
command=php artisan queue:work --queue=validators --tries=3 --daemon
process_name=%(program_name)s_%(process_num)02d
directory=/var/www
autostart=true
autorestart=true
user=root
numprocs=10
stdout_logfile=/var/www/storage/logs/jobs.out.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/www/storage/logs/jobs.err.log
stderr_logfile_maxbytes=10MB

[program:explorer.worker.broadcast]
command=php artisan queue:work --queue=broadcast --tries=3 --daemon
process_name=%(program_name)s_%(process_num)02d
directory=/var/www
autostart=true
autorestart=true
user=root
numprocs=1
stdout_logfile=/var/www/storage/logs/jobs.out.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/www/storage/logs/jobs.err.log
stderr_logfile_maxbytes=10MB

[program:explorer.worker.broadcast_tx]
command=php artisan queue:work --queue=broadcast_tx --tries=3 --daemon
process_name=%(program_name)s_%(process_num)02d
directory=/var/www
autostart=true
autorestart=true
user=root
numprocs=1
stdout_logfile=/var/www/storage/logs/jobs.out.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/www/storage/logs/jobs.err.log
stderr_logfile_maxbytes=10MB

[program:explorer.worker.block_events]
command=php artisan queue:work --queue=block-events --tries=3 --daemon
process_name=%(program_name)s_%(process_num)02d
directory=/var/www
autostart=true
autorestart=true
user=root
numprocs=1
stdout_logfile=/var/www/storage/logs/jobs.out.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/www/storage/logs/jobs.err.log
stderr_logfile_maxbytes=10MB

[program:explorer.worker.main_queue]
command=php artisan queue:work --queue=main --tries=3 --daemon
process_name=%(program_name)s_%(process_num)02d
directory=/var/www
autostart=true
autorestart=true
user=root
numprocs=1
stdout_logfile=/var/www/storage/logs/main_queue.out.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/www/storage/logs/main_queue.err.log
stderr_logfile_maxbytes=10MB

[group:minterexplorer]
programs=explorer.worker.blocks,explorer.worker.nodes,explorer.worker.transactions,explorer.worker.balance,explorer.worker.validators,explorer.worker.broadcast,explorer.worker.broadcast_tx,explorer.worker.block_events,explorer.worker.main_queue